# smf/smartos-user.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<!--
  This Source Code Form is subject to the terms of the Mozilla Public
  License, v. 2.0. If a copy of the MPL was not distributed with this
  file, You can obtain one at http://mozilla.org/MPL/2.0/.

  Copyright (c) 2019 Joyent, Inc.

  This service creates a user per service instance on each boot.  In order for
  the user to be useful, the password must be set or ssh keys must be present.

  When using this as standalone SmartOS (not as a Triton Compute Node),
  /etc/shadow is preserved across reboots. As such, with SmartOS it is possible
  to preserve the password across reboots while keeping its hash secret.  As a
  Triton Compute Node, you must either specify the shadow/pass property in each
  instance or maintain per-user .ssh/authorized_keys.
-->

<service_bundle type="manifest" name="smartos-user">
  <service name="smartdc/smartos-user" type="service" version="1">
    <exec_method type="method" name="start"
      exec="/opt/custom/bin/add-smartos-user %i" timeout_seconds="30" />
    <exec_method type="method" name="stop" exec=":true" timeout_seconds="60" />

    <!--
      Do not worry about dependencies.  Since this lives in zones/opt, there's
      no chance that it can be run too early.
    -->
    <property_group name="startd" type="framework">
      <propval name="duration" type="astring" value="transient" />
      <propval name="ignore_error" type="astring" value="core,signal" />
    </property_group>

    <!-- One instance per user -->
    <instance name="joyent" enabled="true">
      <!--
        This would be better with astring_list but those are only allowed in
        profiles, not manifests
      -->
      <property_group name="profiles" type="application">
        <propval name="admin" type="astring" value="Primary Administrator" />
      </property_group>

      <property_group name="passwd" type="application">
        <propval name="dir" value="/zones/joyent" type="astring" />
        <propval name="comment" value="Joyent Engineering" type="astring" />
      </property_group>

      <property_group name="shadow" type="application">
        <!--
          Uncomment the following to set the password to "8SPPYvsTGxytE7s6".
          Warning! Whatever you put here is world readable.  Use of ssh keys
          instead is advised.
        -->
        <!-- <propval name="pass" type="astring"
          value="$5$ASERtDf8$PT5CgDKf9Vk4mJAZsgAeC.j.E21XHL6Fkprkle3uIAA" /> -->
      </property_group>
    </instance>

    <!--
      Additonal users can be added.  Be sure they each have a dependency on the
      earlier one because the start method is racy.
    -->
    <!-- Delete this line and the close of the comment 10 lines below to enable.
    <instance name="jack" enabled="true">
      <dependency name="joyent" type="service" restart_on="none" grouping="require_all">
        <service_fmri value="svc:/smartdc/smartos-user:joyent" />
      </dependency>
      <property_group name="passwd" type="application">
        <propval name="dir" value="/zones/jack" type="astring" />
        <propval name="comment" value="Mike Gerdts" type="astring" />
      </property_group>
    </instance>
    -->

    <stability value="Evolving" />
    <template>
      <common_name>
        <loctext xml:lang="C">Joyent User Creator</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
